# framework-led-widgets.service
[Unit]
Description=Framework 16 LED Matrix Widgets
After=multi-user.target
Wants=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/framework-led-widgets --config /etc/framework-led-widgets/config.toml
Restart=always
RestartSec=5
User=root
Group=root

StandardOutput=journal
StandardError=journal
SyslogIdentifier=framework-led-widgets

[Install]
WantedBy=multi-user.target
